# speedtest.service
[Unit]
Description=Speedtest job

[Service]
Type=oneshot
ExecStart=speedtest-cli --simple
StandardOutput=file:/tmp/speedtest.out
